# dyndns-update.service
[Unit]
Description=DynDNS Update Service for INWX
After=network.target

[Service]
Type=simple

# Ändere <BENUTZER> zu dem Benutzernamen, unter dem das Skript laufen soll (z.B. root oder www-data)
User=<BENUTZER>

# Ändere <PFAD_ZUM_PROJEKT> zum absoluten Pfad deines Projektverzeichnisses
# z.B. /home/user/dyndns-service
WorkingDirectory=<PFAD_ZUM_PROJEKT>
ExecStart=/usr/bin/node <PFAD_ZUM_PROJEKT>/src/dydns-update.js

# Lädt die Umgebungsvariablen aus der .env-Datei im Projektverzeichnis
EnvironmentFile=<PFAD_ZUM_PROJEKT>/.env

# Die folgende Zeile sorgt dafür, dass der Dienst neu gestartet wird.
# Das Intervall wird aus der .env-Datei gelesen (Standard: 300s).
Restart=on-failure
RestartSec=${EXECUTION_INTERVAL}

[Install]
WantedBy=multi-user.target
